# deploy/discord-bot.service
[Unit]
Description=Discord Bot (Rust)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/discord-bot
EnvironmentFile=/home/pi/discord-bot/.env
ExecStart=/home/pi/discord-bot/target/release/discord-bot
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
